# obunatrans.conf
[program:googletrobot]
directory=/home/ubuntu/GoogleTranslateoBot
command=python3 app.py

autostart=true
autorestart=true
environment=HOME="/home/ubuntu",USER="ubuntu"

stderr_logfile=/home/ubuntu/GoogleTranslateoBot/logfile_err.log
stdout_logfile=/home/ubuntu/GoogleTranslateoBot/logfile.log